# etc/systemd/system/rssgrabber.timer
[Unit]
Description=RSSGrabber timer

[Timer]
OnUnitActiveSec=30min


[Install]
WantedBy=multi-user.target

